# templates/upstart.conf.erb
# Odoo

description  "Odoo"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

setuid <%= @odoo_user %>
setgid <%= @odoo_group %>

exec <%= @install_path %>/odoo.py --config <%= @config_path %>/odoo.conf 